# etc/restic/backup-restic/files.list
/etc/restic
/etc/systemd/system/restic.*
