# etc/systemd/system/restic.timer
[Unit]
Description=Run Restic at 03:10 AM

[Timer]
OnCalendar=*-*-* 03:10:00

[Install]
WantedBy=timers.target
